# docker/php-apache/crontab
* * * * * cd /var/www/html/ && ./vendor/bin/crunz schedule:run
